# docker/files/supervisord.conf
; Sample supervisor config file.

[program:price_monitor]
command=/crypto/bin/price_monitor
directory=/crypto/bin/
user=root
autostart=true
autorestart=true
stdout_logfile=/run_crypto/log/price_monitor/stdout
stderr_logfile=/run_crypto/log/price_monitor/stderr
; the end

[program:account_monitor]
command=/crypto/bin/account_monitor
directory=/crypto/bin/
user=root
autostart=true
autorestart=true
stdout_logfile=/run_crypto/log/account_monitor/stdout
stderr_logfile=/run_crypto/log/account_monitor/stderr
; the end

[program:process_delegator]
command=/crypto/bin/account_process_delegator
directory=/crypto/bin/
user=root
autostart=true
autorestart=true
stdout_logfile=/run_crypto/log/process_delegator/stdout
stderr_logfile=/run_crypto/log/process_delegator/stderr
; the end

[program:message_delegator]
command=/crypto/bin/message_delegator
directory=/crypto/bin/
user=root
autostart=true
autorestart=true
stdout_logfile=/run_crypto/log/message_delegator/stdout
stderr_logfile=/run_crypto/log/message_delegator/stderr
; the end

[unix_http_server]
file=/var/run/supervisor.sock   ; (the path to the socket file)
chmod=0700                       ; sockef file mode (default 0700)

[supervisord]
logfile=/var/log/supervisor/supervisord.log ; (main log file;default $CWD/supervisord.log)
pidfile=/var/run/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
childlogdir=/var/log/supervisor            ; ('AUTO' child log dir, default $TEMP)

; the below section must remain in the config file for RPC
; (supervisorctl/web interface) to work, additional interfaces may be
; added by defining them in separate rpcinterface: sections
[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock ; use a unix:// URL  for a unix socket

; The [include] section can just contain the "files" setting.  This
; setting can list multiple files (separated by whitespace or
; newlines).  It can also contain wildcards.  The filenames are
; interpreted as relative to this file.  Included files *cannot*
; include files themselves.

[include]
files = /etc/supervisor/conf.d/*.conf
